# long_script.conf
[program:long_run_script]
command=/usr/local/bin/long_script.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/long_script_err.log
stdout_logfile=/var/log/long_script_out.log
